# config/crontab
0 * * * * /var/www/html/scripts/speedtest.py >/dev/stdout 2>&1