# config/crontab
0 * * * * /var/www/html/scripts/speedtest.py >/dev/null 2>&1